# debian/cm4all-spawn-accessory.socket
[Socket]
ListenSequentialPacket=@cm4all-spawn
PassCredentials=true
