# debian/cm4all.slice
[Unit]
Description=CM4all Process Spawner Slice
Before=slices.target
